# lantern/systemd/lantern.service
[Unit]
Description=LANtern dashboard server
After=network.target
Wants=network.target

[Service]
Type=simple
User=lantern
Group=lantern
WorkingDirectory=/opt/lantern
Environment=LANTERN_SERVER_CONFIG=/etc/lantern/server.conf
ExecStart=/opt/lantern/venv/bin/python -m server.app.serve
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
RuntimeDirectory=lantern
RuntimeDirectoryMode=0755
ExecStartPre=/usr/bin/mkdir -p /run/lantern/config /run/lantern/snapshots
ExecStartPre=/bin/sh -c '[ -f /run/lantern/config/updateFrequency ] || printf "120\n" > /run/lantern/config/updateFrequency'
ReadOnlyPaths=/opt/lantern /etc/lantern/server.conf /var/lib/success
ReadWritePaths=/var/lib/lantern/state /run/lantern

[Install]
WantedBy=multi-user.target
